# etc/service/electrumx/run
#!/bin/sh
exec 2>&1

mkdir -p /home/electrumx/leveldb

exec python3 /electrumx/electrumx_server
